# gemini/multimodal-live-api/websocket-demo-app/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:websocket]
command=python3 /app/main.py # use python3 command
directory=/app